# bluetooth.service.d/exec.conf
# Drop-in override for bluetooth.service
# Enables BlueZ --compat mode (required for Profile1.RegisterProfile) and
# disables every plugin that would otherwise grab the HID/HOG/A2DP/etc PSMs
# we want to own ourselves.
#
# Install path on device: /etc/systemd/system/bluetooth.service.d/exec.conf
# Reload after edit:      systemctl daemon-reload && systemctl restart bluetooth

[Service]
# Clear the stock ExecStart before defining our own (systemd otherwise stacks them).
ExecStart=
ExecStart=/usr/libexec/bluetooth/bluetoothd --compat --noplugin=deviceinfo,input,hog,battery,scanparam,audio,a2dp,avrcp,hostname,wiimote,autopair,policy,network,opp,pbap,sap,health,midi,adminpolicy,diag,mesh,bap,vcp,mcp,ccp,gatt,media -d -n
StandardOutput=append:/var/log/bluetoothd.log
StandardError=append:/var/log/bluetoothd.log

# Notes:
# - If `bluetoothd` lives at /usr/lib/bluetooth/bluetoothd on your firmware
#   (Debian puts it there instead of /usr/libexec/), edit the path above.
#   Detect at install time with:
#       ls /usr/libexec/bluetooth/bluetoothd /usr/lib/bluetooth/bluetoothd 2>/dev/null
# - The plugin blocklist is intentionally aggressive — anything that might
#   register Profile1 under standard HID/HOG/Audio UUIDs and steal our PSMs.
# - -d -n keeps bluetoothd in foreground with debug logging. Combined with
#   Type=dbus in the stock unit this works because bluetoothd sends the
#   READY=1 signal dbus needs. If your stock unit is Type=notify and fails
#   to start, drop the -n flag.
